# 50-evcape.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", TAG+="systemd", RUN+="/usr/bin/systemctl restart evcape"
