# evcape.timer
[Unit]
Description=Run evcape at boot with a delay

[Timer]
OnBootSec=15s

[Install]
WantedBy=timers.target

